# etc/profile.d/10-parch-env.sh
export EDITOR=vim
export VISUAL=vim
export PAGER=less
export LESS='-R --use-color -Dd+r$Du+b'
export MANPAGER='less -R'
export SYSTEMD_LESS='FRSXMK'
